# docker/supervisor.conf
[supervisord]
nodaemon = true

[unix_http_server]
file = /var/run/supervisor.sock 

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:zabbix_agent]
command = /usr/sbin/zabbix_agent2 -c /etc/zabbix/zabbix_agentd.conf