# apache/supervisor.conf
[supervisord]
nodaemon=true

[program:apache]
command=httpd-foreground

[program:zabbix_agent]
command=/usr/sbin/zabbix_agentd -c /etc/zabbix/zabbix_agentd.conf