# redis/supervisor.replica.conf
[supervisord]
nodaemon=true

[program:redis]
command=redis-server --replicaof redis-master 6379

[program:zabbix_agent]
command=/usr/sbin/zabbix_agent2 -c /etc/zabbix/zabbix_agentd.conf