# redis/supervisor.master.conf
[supervisord]
nodaemon=true

[program:redis]
command=redis-server --appendonly yes

[program:zabbix_agent]
command=/usr/sbin/zabbix_agent2 -c /etc/zabbix/zabbix_agentd.conf