# php-fpm/php-fpm_agent/zabbix_agent_supervisor.conf

[program:zabbix_agent]
command=zabbix_agentd